# systemd/.config/systemd/user/i3-session-save.service
[Unit]
Description=Save i3 session (layout + programs) on logout/shutdown
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=%h/.local/bin/i3-session-save.sh

[Install]
WantedBy=i3-session.target
